# deploy/vitaicare.service
# Servicio systemd para el dashboard VITAICARE en la Raspberry Pi.
#
# Instalación (ver README_RASPBERRY.md):
#   1. Editá User, WorkingDirectory y ExecStart según tu usuario y la ruta del repo.
#   2. sudo cp deploy/vitaicare.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now vitaicare
#
# Usa gunicorn con UN solo worker a propósito: cada worker carga toda la serie
# de datos en memoria, y en una Pi de 2 GB más de un worker provocaría OOM.

[Unit]
Description=VITAICARE - Dashboard de monitoreo de pacientes (Dash)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# CAMBIAR: usuario de la Raspberry (p. ej. 'pi' o el tuyo).
User=pi
# CAMBIAR: ruta absoluta al repositorio en la Raspberry.
WorkingDirectory=/home/pi/newpacientes
# CAMBIAR: ruta al gunicorn dentro del entorno virtual del repo.
ExecStart=/home/pi/newpacientes/.venv/bin/gunicorn \
    --workers 1 \
    --threads 2 \
    --timeout 120 \
    --bind 0.0.0.0:8050 \
    app:server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
